# src/scrp.supervisor.conf
[program:scrp]
command=/bin/bash -c "cd /app/scrp && nice -n 5 sleep 10000000"
redirect_stderr=true
